# rootfs/common/etc/init.d/bcreset
#!/bin/sh /etc/rc.common
#
# OpenSync reference board specifics
#

START=99

start()
{
    /usr/plume/tools/bootconfig --write boot_log | logger -t BCRESET
}
