# rootfs/common/etc/init.d/fan
#!/bin/sh /etc/rc.common
#
# OpenSync reference board specifics
#

START=60

start()
{
    /usr/plume/bin/fan.sh 5500 | logger -t FAN
}
