# rootfs/common/etc/init.d/firewall
#!/bin/sh /etc/rc.common
# {# jinja-parse #}
INSTALL_PREFIX={{INSTALL_PREFIX}}

START=60

start()
{
    ${INSTALL_PREFIX}/bin/iptables_cmd.sh boot
}
