# rootfs/common/etc/init.d/firewall
#!/bin/sh /etc/rc.common

START=60

start()
{
    /usr/plume/bin/iptables_cmd.sh boot
}
