# rootfs/kconfig/PLATFORM_QCA_QSDK110/etc/init.d/dnsmasq
#!/bin/sh /etc/rc.common

#
# Adding a dummy script, since dnsmasq will be started by OpenSync.
# It should not be running by default while device bootup.
#
